# btc-paper-trading.service
[Unit]
Description=BTC Agent System — Paper Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/btc-agent-system
Environment=PATH=/home/ubuntu/.hermes/hermes-agent/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/ubuntu/.hermes/hermes-agent/venv/bin/python main.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Graceful shutdown (send SIGTERM, wait 15s)
TimeoutStopSec=15
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
